# roles/user-stuff/templates/logrotate.cron.j2
#!/bin/sh

test -x /usr/sbin/logrotate || exit 0
/usr/sbin/logrotate ~{{ user }}/.local/etc/logrotate.conf
