# roles/user-stuff/files/logrotate.timer
[Unit]
Description=Trigger logrotate

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
